# conf/rc-p.service.template

# Description:
# Default opensearch.service file

[Unit]
Description=rc-ticket

[Service]
User=www-data
Group=www-data
Environment="APP_RUN_PATH=/var/run/repair-cafe-app/"
ExecStart=/var/www/html/repair-cafe-app/bin/start-printer-p1 --mac {{ MAC }}

[Install]
WantedBy=multi-user.target

